# wiki/config/sv/nginx/finish
#!/bin/sh

# Kill our processes
pkill -TERM runsvdir
