# wiki/rootfs/etc/services.d/php/run
#!/usr/bin/with-contenv sh

/usr/local/sbin/php-fpm --nodaemonize -R
