# wiki/config/sv/php/finish
#!/bin/sh

# kill off the nginx
pkill -TERM runsvdir
