# wiki/config/services.d/php/run
#!/usr/bin/execlineb -P

s6-envdir /usr/local/env /usr/local/sbin/php-fpm --nodaemonize -R
